# deploy/accounting-max.service
[Unit]
Description=Accounting MAX bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Каталог клона: .venv, .env, run.py. Подставь свой путь, если не /root/yandex_tables_parse
WorkingDirectory=/root/yandex_tables_parse
Environment=PYTHONUNBUFFERED=1
# Через /bin/sh: один бинарный путь без споров парсера; при 203 на «прямом» python часто виноват CRLF (^M) в unit-файле.
ExecStart=/bin/sh -c 'exec /root/yandex_tables_parse/.venv/bin/python /root/yandex_tables_parse/run.py'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
