# deploy/accounting-max.service
[Unit]
Description=Accounting MAX bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Каталог клона: .venv, .env, run.py. Подставь свой путь, если не /root/yandex_tables_parse
WorkingDirectory=/root/yandex_tables_parse
Environment=PYTHONUNBUFFERED=1
# Прямой запуск Python — без bash-скрипта (нет systemd 203/EXEC из-за shebang/CRLF/правок run-accounting-max.sh).
# Нужен: python3 -m venv .venv && .venv/bin/pip install -r requirements.txt
ExecStart=/root/yandex_tables_parse/.venv/bin/python /root/yandex_tables_parse/run.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
