# deploy/accounting-max.service
[Unit]
Description=Accounting MAX bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Должен совпадать с каталогом клона (где .venv, .env, run.py). Пример:
WorkingDirectory=/root/yandex_tables_parse
# Явно через bash: иначе ядро гоняет shebang — при CRLF или /usr/bin/env bash → systemd 203/EXEC
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash /root/yandex_tables_parse/deploy/run-accounting-max.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
